# 7.defense/misc.crontab
# Weekly reboot, runs Sunday at 4:30 a.m.
30 4 * * 0 /sbin/reboot

# Time sync, runs every 15 minutes
15 * * * * /usr/sbin/ntpdate time.nist.gov >/dev/null 2>&1
